# sonoff-spm-bridge/rootfs/etc/services.d/spm-bridge/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Sonoff SPM Bridge stopped."
