# sonoff-spm-bridge/rootfs/etc/services.d/spm-bridge/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Sonoff SPM Bridge..."
exec python3 /app/main.py
